# config/systemd/user/wallpaper-rotate.service
[Unit]
Description=Rotate sway wallpaper

[Service]
Type=oneshot
ExecStart=%h/.config/sway/wallpaper-rotate.sh
